# deploy/smart-home-monitor.service
[Unit]
Description=Smart Home Monitor Server
After=network-online.target mysql.service remote-fs.target
Wants=network-online.target

[Service]
Type=simple
User=smart-home
Group=smart-home
WorkingDirectory=/opt/smart-home
ExecStart=/opt/smart-home/bin/SmartHomeMonitorServer /etc/smart-home/server.conf
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/smart-home /mnt/smart-home-recordings
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
